# systemd/obsidian-daily-note.timer
[Unit]
Description=Run a script every day

[Timer]
OnCalendar=*-*-* 03:05:00
Persistent=true

[Install]
WantedBy=timers.target
