# systemd/obsidian-sync.service
[Unit]
Description=Obsidian vault sync (continuous)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sam
Environment=PATH=/home/sam/.nvm/versions/node/v22.22.1/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/sam/.nvm/versions/node/v22.22.1/bin/ob sync --continuous --path /home/sam/obsidian/Vault
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
