# systemd/update-dns.timer
[Unit]
Description=Run update-dns every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
